# infra/systemd/amber-tmux.service
[Unit]
Description=amber-ide tmux server (persistence spine)

[Service]
Type=oneshot
RemainAfterExit=yes
# Start the tmux server on the default socket if it is not already running.
# Sourcing tmux.conf on server start triggers tmux-continuum's restore, which
# recreates all saved sessions; "_amber-boot" is a throwaway bootstrap session
# that keeps the server alive while the restore runs (a post-restore hook in
# tmux.conf removes it).
ExecStart=/bin/sh -c 'tmux has-session 2>/dev/null || exec tmux new-session -d -s _amber-boot'

[Install]
WantedBy=default.target
