# infra/launchd/com.amber-ide.tmux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amber-ide.tmux</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>PATH="/opt/homebrew/bin:/usr/local/bin:$PATH"; tmux has-session 2>/dev/null || exec tmux new-session -d -s _amber-boot "$HOME/.local/share/amber-ide/bin/boot-hold.sh"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
